# deploy_tools/supervisor_recognize_service_beat.conf
; ================================
;  celery beat supervisor example
; ================================
; the name of your supervisord program
[program:recognize_service_beat]

; Set full path to celery program if using virtualenv
command=RECOGNIZE_SERVICE_DIRECTORY/venv/bin/celery -A service.conqueror.scheduling beat -l info

; The directory to your Django project
directory=RECOGNIZE_SERVICE_DIRECTORY

; If supervisord is run as the root user, switch users to this UNIX user account before doing any processing.
user=root

; Supervisor will start as many instances of this program as named by numprocs
numprocs=1

; Put process stdout output in this file
stdout_logfile=/var/log/celery/recognize_service_beat.log
; Put process stderr output in this file
stderr_logfile=/var/log/celery/recognize_service_beat.log
; If true, this program will start automatically when supervisord is started
autostart=true
; May be one of false, unexpected, or true. If false, the process will never be autorestarted. If unexpected, the process will be resta$
autorestart=true
; The total number of seconds which the program needs to stay running after a startup to consider the start successful.
startsecs=10
; if your broker is supervised, set its priority higher so it starts first
priority=999
